# server/monthly-refresh.cron
# Run at 09:30 on the 5th day of every month.
# Replace /opt/feishu-online-sheets-skill with the actual server path.
30 9 5 * * cd /opt/feishu-online-sheets-skill && /bin/bash server/run_refresh.sh
