# server/weekly-refresh.cron
# Run every Monday at 09:30 and refresh through the latest completed Sunday.
# Replace /opt/feishu-project-delivery-cycle-updater with the actual server path.
30 9 * * 1 cd /opt/feishu-project-delivery-cycle-updater && PYTHON_BIN=/opt/feishu-project-delivery-cycle-updater/.venv/bin/python BITABLE_PUBLISHER=/opt/feishu-project-delivery-cycle-updater/feishu-online-sheets/scripts/publish_bitable.py /bin/bash server/run_refresh.sh
