# server/weekly-refresh.cron
# Run every Monday at 08:00 and refresh the latest current-quarter snapshot through the latest completed Sunday.
# Replace /opt/feishu-project-delivery-intelligence with the actual server path.
0 8 * * 1 cd /opt/feishu-project-delivery-intelligence && PYTHON_BIN=/opt/feishu-project-delivery-intelligence/.venv/bin/python BITABLE_PUBLISHER=/opt/feishu-project-delivery-intelligence/feishu-online-sheets/scripts/publish_bitable.py /bin/bash server/run_refresh.sh
